# server files/supervisor_wey.conf
[program:wey]
command = /webapps/wey/env/bin/gunicorn_start
user = weyuser
stdout_logfile = /webapps/wey/logs/supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8